# docker/spider/crontab
0 6 * * * go run /go/src/spider/main.go -mode 0